# package/debian1/DEBIAN/prerm
#!/bin/bash
systemctl stop tellod-ad.service
systemctl disable tello-ad.service
